# src/bios-agent-rt.service.in
[Unit]
Description=BIOS have info about all current metrics and is able to share it via REQ-REP
After=malamute.service
PartOf=bios.target

[Service]
Type=simple
Environment="prefix=@prefix@"
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
ExecStart=@prefix@/bin/bios-agent-rt

[Install]
WantedBy=bios.target
