# src/fty-metric-cache.service.in
[Unit]
Description=fty-metric-cache service
Requires=malamute.service
After=malamute.service
PartOf=bios.target

[Service]
Type=simple
User=bios
Environment="prefix=@prefix@"
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
ExecStart=@prefix@/bin/fty-metric-cache
Restart=always

[Install]
WantedBy=bios.target


